# src/soyocontroller.service.in
[Unit]
Description=SoyoSource inverter controller
After=network.target

[Service]
User=nobody
Group=nogroup
ExecStart=@bindir@/soyocontroller
EnvironmentFile=@sysconfdir@/default/soyocontroller

[Install]
WantedBy=multi-user.target
